# DEBIAN/postinst
systemctl enable inventorius-landing-page
systemctl start inventorius-landing-page