# DEBIAN/prerm
systemctl stop inventorius-landing-page
systemctl disable inventorius-landing-page